# layout/DEBIAN/prerm
#!/bin/sh
declare -a cydia
cydia=($CYDIA)

plutil -key EnvironmentVariables -remove /System/Library/LaunchDaemons/com.apple.CommCenterClassic.plist > /dev/null

cd /Applications/Preferences.app/
MESSAGES_PLIST=/System/Library/PrivateFrameworks/Conference.framework/Messages.plist
ENTRY=16

echo "Remove our settings"
plutil -key items -key $ENTRY -key cell -remove $MESSAGES_PLIST > /dev/null
plutil -key items -key $((ENTRY+1)) -key bundle -remove $MESSAGES_PLIST > /dev/null
plutil -key items -key $((ENTRY+1)) -key cell -remove $MESSAGES_PLIST > /dev/null
#plutil -key items -key $((ENTRY+1)) -key hasIcon -remove $MESSAGES_PLIST > /dev/null
plutil -key items -key $((ENTRY+1)) -key isController -remove $MESSAGES_PLIST > /dev/null
plutil -key items -key $((ENTRY+1)) -key label -remove $MESSAGES_PLIST > /dev/null

echo "Restart CommCenter"
launchctl unload /System/Library/LaunchDaemons/com.apple.CommCenterClassic.plist
launchctl load /System/Library/LaunchDaemons/com.apple.CommCenterClassic.plist

if [[ ${CYDIA+@} ]]; then
eval "echo 'finish:restart' >& ${cydia[0]}"
fi
